# scripts/etc/launchd/com.k3d-manager.webhook.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.k3d-manager.webhook</string>
  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/python3</string>
    <string>{{WEBHOOK_BIN}}</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>K3DM_REPO_ROOT</key>
    <string>{{K3DM_REPO_ROOT}}</string>
    <key>KUBECONFIG</key>
    <string>{{HOME}}/.kube/config:{{HOME}}/.kube/ubuntu-k3s.yaml:{{HOME}}/.kube/k3s-gcp.yaml</string>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>{{WEBHOOK_LOG}}</string>
  <key>StandardErrorPath</key>
  <string>{{WEBHOOK_LOG}}</string>
</dict>
</plist>
